# deploy/systemd/hl-api-snapshot-account-values.timer
# Daily per-user account-value snapshot (powers the 24h change). 00:05 UTC.
# No-op until users have wallets on their profile (post-Privy).
[Unit]
Description=hl-api: daily account-value snapshot

[Timer]
OnCalendar=*-*-* 00:05:00 UTC
Persistent=true
Unit=hl-api-job@snapshot-account-values.service

[Install]
WantedBy=timers.target
