# deploy/systemd/hl-api-sync-assets.timer
# Tier A (discovery): full asset catalog sync. Daily.
[Unit]
Description=hl-api: daily asset catalog sync

[Timer]
OnCalendar=*-*-* 08:00:00 UTC
# Run shortly after boot too, so a box that was down comes back up in sync.
OnBootSec=1min
# And catch up a run missed while the box was down.
Persistent=true
Unit=hl-api-job@sync-assets.service

[Install]
WantedBy=timers.target
